# scheduling/my-diary.service
[Unit]
Description=my-diary — Automatyczny skryba dzienny
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
TimeoutStartSec=600
WorkingDirectory=__PROJECT_DIR__
ExecStart=__PROJECT_DIR__/.venv/bin/python -m my_diary
Environment=PATH=%h/.local/bin:/usr/local/bin:/usr/bin:/bin

[Install]
WantedBy=default.target
